# packer/ansible/docker-up.service
[Unit]
Description=Up Docker-compose

[Service]
ExecStart=/usr/bin/docker-up.sh

[Install]
WantedBy=multi-user.target